# doc/flamejam.cron.d
* * * * * root /usr/bin/curl http://domain.tld/tick > /dev/null 2>&1
