# doc/flamejam.cron.d
* * * * * root /usr/bin/curl http://domain.tld/tick
